# deployment/services/codeui.service
[Unit]
Description=CodeUI - AI CLI Tools Management Application
After=network.target
Documentation=https://github.com/managedcode/CodeUI

[Service]
Type=notify
ExecStart=/opt/codeui/CodeUI.Web --urls http://localhost:5000
User=codeui
Group=codeui
WorkingDirectory=/opt/codeui
Restart=always
RestartSec=5
TimeoutStopSec=20
SyslogIdentifier=codeui
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=DOTNET_ENVIRONMENT=Production

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/codeui/Data /opt/codeui/Logs

[Install]
WantedBy=multi-user.target